# bin/tiny-screen-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alrocar.tiny-screen-monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/tiny-screen-monitor</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/opt/homebrew/var/log/tiny-screen-monitor/error.log</string>
    <key>StandardOutPath</key>
    <string>/opt/homebrew/var/log/tiny-screen-monitor/output.log</string>
</dict>
</plist> 